# conf/cam.conf
[program:cam]
command=python /home/pi/flask/cameraPi/cam.py
autostart=false
autorestart=true
stdout_logfile=/home/pi/flask/logs/cam.log
stderr_logfile=/home/pi/flask/logs/cam.log
