# conf/cam.conf
[program:cam]
command=python /home/pi/client/cameraPi/cam.py
autostart=false
autorestart=true
stdout_logfile=/var/log/camserv/camserv.log
stderr_logfile=/var/log/camserv/camserv.log
